# comms/cdrd/files/cdrd.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: cdrd
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable cdrd:
# cdrd_enable (bool):      Set it to "YES" to enable cdrd.
#                               Default is "NO".
# cdrd_conf_file (file):   Set full path to configuration file.
#                               Default is "/usr/local/etc/cdrd/cdrd.conf".
# cdrd_pid_file (path):    Set full path to pid file.
#                               Default is "/var/run/cdrd/cdrd.pid".

. /etc/rc.subr

name="cdrd"
rcvar=${name}_enable

load_rc_config $name

: ${cdrd_enable:="NO"}
: ${cdrd_conf_file:="/usr/local/etc/${name}/${name}.conf"}
: ${cdrd_pid_file:="/var/run/${name}/${name}.pid"}

: ${cdrd_user="uucp"}
: ${cdrd_group="dialer"}

command=/usr/local/bin/${name}
cdrd_flags="-c ${cdrd_conf_file} -n ${name} -p ${cdrd_pid_file}"
command_args="&"
command_interpreter="/bin/sh"

pidfile=${cdrd_pid_file}

run_rc_command "$1"
